# recipes-core/kiosk-refresh/files/kiosk-reload.service
[Unit]
Description=Run reload in chromium browser

[Service]
Type=oneshot
ExecStart=/usr/bin/refresh-chromium.sh
User=kiosk
