# recipes-graphics/x11-common/files/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager
Wants=network.target network-online.target time-sync.target
After=network.target network-online.target time-sync.target

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=always

[Install]
Alias=display-manager.service
